# backup/backup-home.service
/etc/systemd/system/backup-home.service
[Unit]
Description=Backup home directory with Borg

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-home.sh
